# packaging/tpfan-daemon.service
[Unit]
Description=tpfan — ThinkPad Fan Control Daemon
After=systemd-modules-load.service

[Service]
Type=dbus
BusName=org.tpfan1
ExecStart=/usr/local/bin/tpfan-daemon
WatchdogSec=10
Restart=on-failure
NotifyAccess=main
Environment=TPFAN_LOG=info

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/etc/tpfan /var/lib/tpfan /proc/acpi/ibm/fan
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target
